# moodle/crontab
# Moodle Cron 
* * * * * /usr/local/bin/php -c /usr/local/etc/php/php.ini /var/www/html/admin/cli/cron.php >> /var/log/cron.log 2>&1
# End Cron
